# supervisord.conf
[supervisord]
nodaemon=true

[program:factomd]
command=factomd
user=factom

[program:fctwallet]
command=fctwallet
user=factom
redirect_stderr=true
stdout_logfile=/home/factom/.factom/fctwallet.log
