# supervisord.conf
[supervisord]
nodaemon=true

[program:factomd]
command=factomd

[program:fctwallet]
command=fctwallet
redirect_stderr=true
stdout_logfile=/home/factom/.factom/fctwallet.log
