# supervisord.conf
[supervisord]
nodaemon=true

[program:factomd]
command=factomd
user=factom

[program:fctwallet]
command=fctwallet
user=factom
