# supervisord.conf
[supervisord]
nodaemon=true

[program:factomd]
command=factomd
user=factom

[program:factom-walletd]
command=factom-walletd
user=factom
